# cardano-bot.service
[Unit]
Description=Cardano DeFi HFT Exhaustion Bot
Documentation=https://github.com/your-repo/cardano-bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agile
WorkingDirectory=/home/agile/exhaustion-bot
# Use unbuffered output for real-time logging
Environment=PYTHONUNBUFFERED=1
# Use uv for execution
ExecStart=/home/agile/.local/bin/uv run /home/agile/exhaustion-bot/paper_trader.py

# Restart policy for high availability
Restart=always
RestartSec=10s
StartLimitInterval=300
StartLimitBurst=5

# Security hardening
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=append:/home/agile/exhaustion-bot/bot_service.log
StandardError=append:/home/agile/exhaustion-bot/bot_service_error.log

[Install]
WantedBy=multi-user.target
